# 1.8.xx/local/luci/luci-0.11.1/themes/custom/root/etc/uci-defaults/luci-theme-custom
uci batch <<-EOF
        set luci.themes.Custom=/luci-static/custom
        commit luci
EOF
